# deploy/glak-api.service
[Unit]
Description=Genesis Learning Academy API (contact + enrollment forms)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/glak/htdocs/genesislearningacademyofkennesaw.com/public
Environment=NODE_ENV=production
Environment=PORT=3002
EnvironmentFile=-/etc/glak-api.env
ExecStart=/home/glak/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
